# .VARIANT/xyzy/System/ssh/sshd.socket
[Unit]
Conflicts=sshd.service
Wants=sshdgenkeys.service

[Socket]
ListenStream=25800
Accept=yes

[Install]
WantedBy=sockets.target
