# config/systemd/user/a2ln.service
[Unit]
Description=A2LN server

[Service]
Type=simple
ExecStart=/usr/local/bin/a2ln 8888

[Install]
WantedBy=default.target
