# .VARIANT/rin/System/ssh/sshd.socket
[Unit]
Conflicts=sshd.service
Wants=sshdgenkeys.service

[Socket]
ListenStream=25801
Accept=yes

[Install]
WantedBy=sockets.target
