# system/etc/systemd/system/enable-pci-rescan.service
[Unit]
Description=Enable PCI rescan

[Service]
Type=oneshot
ExecStart=/usr/local/bin/enable-pci-rescan

[Install]
WantedBy=multi-user.target
